# server/setup/05-service-settings/etc/profile.d/profile.local.sh
# /etc/profile.d/profile.local.sh

test -e /opt-linux-x86_64/etc/profile.local    && . /opt-linux-x86_64/etc/profile.local

